# odessa/evil-timer/python2-certbot.timer
[Unit]
Description=getting real evil

[Timer]
OnUnitActiveSec=30s

[Install]
WantedBy=timers.target
