# odessa/evil-timer/python2-certbot.service
[Unit]
Description=Drink water reminder

[Service]
Type=oneshot
ExecStart=/usr/bin/echo "get evil NOW"
